# dist/etc/init.d/xrdp
#!/bin/sh
### BEGIN INIT INFO
# Provides:          xrdp
# Required-Start:    $remote_fs $syslog $network
# Required-Stop:     $remote_fs $syslog $network
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: xrdp (X Remote Desktop Protocol) full server
# Description:       The xrdp daemon uses the Remote Desktop Protocol
#                    to present a graphical login to a remote client,
#                    allowing connections to an xorgxrdp or VNC server
#                    or another RDP server. The xrdp-sesman daemon
#                    authenticates the users against PAM and starts
#                    the session and, if necessary, X11 server.
### END INIT INFO
#-
# Copyright © 2015, 2019 mirabilos <thorsten.glaser@teckids.org>
# Published under The MirOS Licence.

# absolute basics
LC_ALL=C PATH=/sbin:/usr/sbin:/bin:/usr/bin
export LC_ALL PATH
unset LANGUAGE

# exit cleanly if disabled or not installed
test -x /usr/sbin/xrdp || exit 0

# Debian/LSB init script foobar
DESC='Remote Desktop Protocol server'
NAME=xrdp
. /lib/init/vars.sh
test -t 0 && VERBOSE=yes
. /lib/lsb/init-functions

# read options
SESMAN_START=yes
SESMAN_OPTIONS=
XRDP_OPTIONS=
test -r /etc/default/xrdp && . /etc/default/xrdp

# prepare for actions
case $1 in
(status)
        # nothing to do here
        ;;
(start|stop|force-reload|restart|try-restart)
        # check for root; create run-time directories
        . /usr/share/xrdp/socksetup
        ;;
(*)
        # syntax error
        echo >&2 "Usage: $0 {start|stop|status|restart|try-restart|force-reload}"
        exit 3
        ;;
esac

# take action
rv=0
case $1 in
(start)
        echo
        systemctl start xrdp xrdp-sesman
        systemctl status xrdp xrdp-sesman
        ;;
(stop)
        systemctl stop xrdp xrdp-sesman
        ;;
(status)
        systemctl status xrdp xrdp-sesman
        ;;
(force-reload|restart)
        systemctl restart xrdp xrdp-sesman
        ;;
esac
