# dist/etc/skel/.config/autostart/lxqt-compton.desktop
[Desktop Entry]
Exec=compton --vsync=opengl-mswc --glx-use-copysubbuffermesa
X-XFCE-Autostart-Override=true
Hidden=false
