# dist/etc/skel/.config/autostart/WSL PulseAudio.desktop
[Desktop Entry]
Comment=PulseAudio Daemon for WSL
Exec=/usr/bin/pulseaudio --enable-memfd=FALSE --disable-shm=TRUE --high-priority=FALSE --realtime=FALSE --log-target=syslog
GenericName=
Icon=vvave
MimeType=
Name=WSL PulseAudio
Path=
StartupNotify=false
Terminal=false
TerminalOptions=
Type=Application
X-DBUS-ServiceName=
X-DBUS-StartupType=
X-KDE-SubstituteUID=false
X-KDE-Username=
